# deploy/systemd/hermes-agent.service
[Unit]
Description=Hermes Agent Service
Documentation=https://github.com/your-org/hermes-agents
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/hermes-agents/agent-service
ExecStart=/usr/bin/node /home/ubuntu/hermes-agents/agent-service/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-agent

# Environment
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOST=0.0.0.0
Environment=HERMES_LOG_LEVEL=info

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/hermes-agents/agent-service/data
ReadWritePaths=/home/ubuntu/hermes-agents/agent-service/agents

[Install]
WantedBy=multi-user.target